# container/cryofall/gsa-cryofall.container
# SPDX-License-Identifier: AGPL-3.0-or-later
# Copyright (c) 2026 Jonathan D.A. Jewell (hyperpolymath) <j.d.a.jewell@open.ac.uk>
#
# Podman Quadlet unit — CryoFall dedicated server (managed by GSA)
#
# Verpex VPS: metadatastician.art (209.42.26.106), AlmaLinux 9.7, Podman 5.6.0
#
# Install (as the service user on Verpex):
#   cp container/cryofall/gsa-cryofall.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user start gsa-cryofall
#
# Manage via GSA (preferred):
#   GSA Gossamer GUI → CryoFall → Actions panel (start/stop/restart/update/backup)
#   gsa probe 209.42.26.106 6000    ← fingerprint the running server
#
# Manual Podman fallback:
#   podman start cryofall
#   podman stop -t 30 cryofall
#   podman logs --tail 100 cryofall
#
# DNS note:
#   The game server is reachable at cryofall.jewell.nexus:6000 (UDP).
#   The DNS record MUST be grey-cloud (Cloudflare DNS-only, not proxied)
#   because Cloudflare cannot proxy UDP traffic.
#
# Build the image before first start:
#   podman build -t localhost/cryofall-server:latest \
#     -f container/cryofall/Containerfile .

[Unit]
Description=CryoFall dedicated server — managed by Game Server Admin
After=network-online.target gsa-verisimdb.service
Wants=network-online.target
# GSA VeriSimDB must be up to track server config state and drift
Requires=gsa-verisimdb.service

[Container]
Image=localhost/cryofall-server:latest
ContainerName=cryofall
AutoUpdate=local

# ── UDP ports (game traffic + discovery) ─────────────────────────────────
# Must be exposed to the internet.  DNS: cryofall.jewell.nexus grey-cloud.
# Cloudflare cannot proxy UDP — use DNS-only (grey cloud) in Cloudflare panel.
PublishPort=0.0.0.0:6000:6000/udp
PublishPort=0.0.0.0:6001:6001/udp

# ── Persistent volumes ────────────────────────────────────────────────────
# cryofall-game-data   ~3 GB SteamCMD game installation (persists across updates)
# cryofall-world-data  Settings.xml + world saves (managed by GSA config editor)
# cryofall-backups     Timestamped world archives (GSA `backup` action)
Volume=cryofall-game-data:/opt/cryofall:Z
Volume=cryofall-world-data:/data/cryofall:Z
Volume=cryofall-backups:/data/backups:Z

# ── Environment ───────────────────────────────────────────────────────────
Environment=CRYOFALL_INSTALL_DIR=/opt/cryofall
Environment=CRYOFALL_DATA_DIR=/data/cryofall
Environment=CRYOFALL_BACKUP_DIR=/data/backups
Environment=CRYOFALL_APPID=1200170

# ── Health check ──────────────────────────────────────────────────────────
# Verify UDP port 6000 is bound.  start-period=180s allows time for the
# first-run SteamCMD installation to complete before health checks begin.
HealthCmd=ss -ulnp | grep -q ':6000 ' || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=180s

# ── Resource limits ───────────────────────────────────────────────────────
# Verpex VPS has 3.6 GB RAM shared between:
#   Tendenci (PostgreSQL + Redis + Caddy) ~800 MB
#   GSA VeriSimDB instances (8090 + 8091) ~512 MB each
#   CryoFall server: allocated 1.5 GB — adequate for up to 10 players
# Adjust if the VPS load profile changes.
PodmanArgs=--memory=1536m --cpus=1.5

# ── Security ──────────────────────────────────────────────────────────────
# SteamCMD and the game server write to volume mounts — rootfs must not be
# read-only.  Drop all capabilities except those needed for network binding.
NoNewPrivileges=true

# ── Labels ────────────────────────────────────────────────────────────────
Label=dev.stapeln.app=game-server-admin
Label=dev.stapeln.role=cryofall-server
Label=dev.stapeln.clade=gm
Label=dev.stapeln.managed-by=gsa

[Service]
Restart=on-failure
RestartSec=10s
# Allow up to 5 min for first-run SteamCMD installation
TimeoutStartSec=300s
TimeoutStopSec=60s

[Install]
WantedBy=default.target
